# ansible/templates/etc#systemd#system#gps.service
[Unit]
Description=GPS data logger
After=multi-user.target

[Service]
Type=simple
Restart=always
WorkingDirectory=/home/pi/oss_marine_datalogger/collectors
ExecStart=/usr/bin/python3 /home/pi/oss_marine_datalogger/collectors/gps.py

[Install]
WantedBy=multi-user.target

